# provisioning/roles/deploy-project/files/supervisor.conf
[program:uwsgi-sara]
command=uwsgi /home/vagrant/sara/uwsgi.ini
autostart=true
autorestart=true
stdout_logfile=/var/log/uwsgi/sara.log
redirect_stderr=true
exitcodes=0